# appengine/gce-backend/agent/machine-provider-agent.service.tmpl
[Unit]
Description=Machine Provider agent
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=%(agent)s

[Install]
WantedBy=multi-user.target
